# daxctl/90-daxctl-device.rules
ACTION=="add", SUBSYSTEM=="dax",\
  RUN+="/usr/bin/daxctl -C $env{DEVNAME}"
